# deploy/umbrella.service
[Unit]
Description=Umbrella App (Next.js)
After=network.target

[Service]
Type=simple
# The Linux user that owns /opt/umbrella-app (created in the guide).
User=umbrella
WorkingDirectory=/opt/umbrella-app
# Reads DATABASE_URL and JWT_SECRET. Keep this file 0600 and owned by the umbrella user.
EnvironmentFile=/opt/umbrella-app/.env
Environment=NODE_ENV=production
Environment=PORT=3000
# Adjust the node path if you installed Node differently (run `which node`).
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
